# sv/log/run
#!/bin/sh
#
#

set -u -e
dir=/progs/logs
sv_name=surferhearts

mkdir -p         $dir/$sv_name
exec svlogd -ttt $dir/$sv_name
